# public/game-data/1/persistence-options.json
{
  "key": "persistence-options",
  "title": "Persistence Options",
  "description": "Choose a method to maintain persistence on the server.",
  "choices": [
    {
      "label": "Registry Run Keys (T1060)",
      "mitreId": "T1060",
      "hourCost": 1,
      "detectionChance": 30,
      "successChance": 70,
      "next": "persistence-options",
      "requiredItems": [],
      "rewards": ["Persistence"]
    },
    {
      "label": "Startup Items (T1547.001)",
      "mitreId": "T1547.001",
      "hourCost": 1,
      "detectionChance": 20,
      "successChance": 80,
      "next": "persistence-options",
      "requiredItems": [],
      "rewards": ["Persistence"]
    },
    {
      "label": "Scheduled Task (T1053.005)",
      "mitreId": "T1053.005",
      "hourCost": 1,
      "detectionChance": 25,
      "successChance": 75,
      "next": "persistence-options",
      "requiredItems": [],
      "rewards": ["Persistence"]
    },
    {
      "label": "Return to initial foothold",
      "mitreId": "N/A",
      "hourCost": 0,
      "detectionChance": 0,
      "successChance": 100,
      "next": "foothold-workstation",
      "requiredItems": []
    }
  ]
}
